# dynamic-distro/scatest/code-from-elsewhere/bad-systemd/bad.service
[Unit]
FooOption=bla
Description=Foooooo

[Service]
Type=noneshot
ExecStart=/usr/sbin/doesreallynotexist